# systemd/aurora-web.service
[Unit]
Description=Listener for distributing aurora repository

[Service]
ExecStart=/usr/bin/aurorad -L
Restart=always

[Install]
WantedBy=multi-user.target
